# com.user.timemachine-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.timemachine-backup</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-c</string>
		<string>HOME=/Users/josh exec /Users/josh/bin/timemachine-backup-launch.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>7200</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/josh/Library/Logs/timemachine-backup-stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/josh/Library/Logs/timemachine-backup-stdout.log</string>
</dict>
</plist>
